# deploy/firecrackerlacker.service
[Unit]
Description=firecrackerlacker - agentic sandbox platform
After=network.target
Requires=network.target

[Service]
Type=simple
ExecStart=/opt/firecrackerlacker/bin/firecrackerlacker-server \
  --listen :8080 \
  --db /opt/firecrackerlacker/db/firecrackerlacker.db \
  --state-dir /opt/firecrackerlacker/vms \
  --firecracker /usr/local/bin/firecracker \
  --kernel /opt/firecrackerlacker/images/default/vmlinux \
  --rootfs /opt/firecrackerlacker/images/default/rootfs.ext4 \
  --bridge fcbr0 \
  --max-sandboxes 50
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
WorkingDirectory=/opt/firecrackerlacker

[Install]
WantedBy=multi-user.target
