# deploy/firecrackerlacker.service
[Unit]
Description=firecrackerlacker - agentic sandbox platform
After=network.target
Requires=network.target

[Service]
Type=simple
ExecStart=/opt/firecrackerlacker/bin/firecrackerlacker-server \
  --listen :8080 \
  --db /opt/firecrackerlacker/db/firecrackerlacker.db \
  --state-dir /opt/firecrackerlacker/vms \
  --images-dir /opt/firecrackerlacker/images \
  --firecracker /usr/local/bin/firecracker \
  --kernel /opt/firecrackerlacker/images/vmlinux \
  --rootfs /opt/firecrackerlacker/images/default/rootfs.ext4 \
  --bridge fcbr0 \
  --max-sandboxes 50 \
  --max-per-key 10 \
  --rate-limit 30 \
  --max-vcpu 4 \
  --max-mem-mib 2048 \
  --default-vcpu 1 \
  --default-mem-mib 256 \
  --pool-size 1 \
  --snapshot-dir /opt/firecrackerlacker/snapshots \
  --prometheus
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
WorkingDirectory=/opt/firecrackerlacker

[Install]
WantedBy=multi-user.target
